# src/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:pcscd]
command=/usr/sbin/pcscd -a -d -f
autostart=true

[program:vnc]
command=/bin/bash -c "/dockerstartup/vnc_startup.sh --wait"
autostart=true
user=1000


